# tellsticklive/rootfs/etc/services.d/runonce/run
#!/command/with-contenv bashio

s6-svc -O /run/s6/legacy-services/runonce

if bashio::config.true "enable_live" && bashio::config.is_empty "live_uuid"; then
    # Wait for telldusd sockets to be ready before starting tellive
    bashio::log.info "Waiting for telldusd to be ready for registration..."
    SOCKET_WAIT=0
    while [[ ! -S /tmp/TelldusClient ]] || [[ ! -S /tmp/TelldusEvents ]]; do
        sleep 1
        SOCKET_WAIT=$((SOCKET_WAIT + 1))
        if [[ ${SOCKET_WAIT} -ge 60 ]]; then
            bashio::log.error "Timeout waiting for telldusd sockets for registration"
            exit 1
        fi
    done
    
    bashio::log.info "UUID has not been entered in config. Generating registration URL"
    bashio::log.info "Follow the link below to register this addon with Telldus Live."
    bashio::log.info "Copy the uuid string in the url and set live_uuid config option."
    exec tellive_core_connector /etc/tellive.conf
fi
