# tellsticklive/rootfs/etc/s6-overlay/s6-rc.d/tellivecore/finish
#!/command/with-contenv bashio
# ==============================================================================
# Take down the S6 supervision tree when tellive core fails
# ==============================================================================

# Exit code 0 means the service exited normally (e.g., not enabled)
# In this case, exit 125 to prevent restart
if [[ "${1}" -eq 0 ]]; then
    exit 125
fi

# Exit code 256 is expected for s6-svc -Od .
if [[ "${1}" -eq 256 ]]; then
    exit 125
fi

bashio::log.warning "tellive crashed with exit code ${1}, halting add-on"
/run/s6/basedir/bin/halt
