# tellsticklive/rootfs/etc/s6-overlay/s6-rc.d/tellive_setup/up
/etc/s6-overlay/scripts/tellive_setup/tellive_setup.sh