# tellsticklive/rootfs/etc/services.d/stdin/finish
#!/command/with-contenv bashio
# ==============================================================================
# Handle STDIN service finish
# ==============================================================================

# Exit code 0 means the service exited normally
# In this case, exit 125 to prevent immediate restart
if [[ "${1}" -eq 0 ]]; then
    exit 125
fi

bashio::log.info "STDIN service stopped with exit code ${1}, restarting..."
