# tellsticklive/rootfs/etc/s6-overlay/s6-rc.d/crond/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant App: Tellsticklive
# Runs the crond deamon
# ==============================================================================

bashio::log.info "Starting Crond...."
exec /usr/sbin/crond -f