# tellsticklive/rootfs/etc/services.d/tellivecore/run
#!/command/with-contenv bashio

if bashio::config.true "enable_live" && bashio::config.has_value "live_uuid"; then
    bashio::log.info "Telldus Live enabled."
    
    # Wait for telldusd sockets to be ready before starting tellive
    bashio::log.info "Waiting for telldusd to be ready..."
    SOCKET_WAIT=0
    while [[ ! -S /tmp/TelldusClient ]] || [[ ! -S /tmp/TelldusEvents ]]; do
        sleep 1
        SOCKET_WAIT=$((SOCKET_WAIT + 1))
        if [[ ${SOCKET_WAIT} -ge 60 ]]; then
            bashio::log.error "Timeout waiting for telldusd sockets"
            exit 1
        fi
    done
    bashio::log.info "Telldusd is ready"
    
    # Additional delay as configured (for sensor discovery, etc.)
    LIVE_DELAY=$(bashio::config 'live_delay')
    if [[ ${LIVE_DELAY} -gt 0 ]]; then
        bashio::log.info "Waiting ${LIVE_DELAY} seconds before connecting to Telldus Live..."
        sleep "${LIVE_DELAY}"
    fi
    
    bashio::log.info "UUID has been entered in config options. Starting Telldus Live..."
    exec tellive_core_connector /etc/tellive.conf
fi

# Check if live is enabled but UUID is not yet configured (registration phase)
if bashio::config.true "enable_live" && bashio::config.is_empty "live_uuid"; then
    bashio::log.info "Telldus Live: Waiting for registration to complete via runonce service."
    bashio::log.info "The tellivecore service will remain disabled until live_uuid is configured."
else
    bashio::log.info "Telldus Live is not enabled. Service disabled."
fi

exec s6-svc -Od .